# services/scan-normalizer.service
[Unit]
Description=TurtleBot4 Scan Normalizer
After=network.target turtlebot4.service
Wants=turtlebot4.service

[Service]
Type=simple
User=ubuntu
ExecStart=/bin/bash -c "source /etc/turtlebot4/setup.bash && python3 /opt/turtlebot4/scripts/scan_normalizer.py"
Restart=on-failure
RestartSec=5
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
